# demonize/panon_connectto_pulseeffects.service
[Unit]
Description=Always Activate PulseEffects

[Service]
Type=simple
StandardOutput=null
ExecStart=

[Install]
WantedBy=default.target
